# examples/systemd/jellyfin_exporter.socket
[Unit]
Description=Jellyfin Exporter

[Socket]
ListenStream=9594

[Install]
WantedBy=sockets.target
